# backstage-portal/deploy/backstage-portal.service
[Unit]
Description=BaselithCore Backstage developer portal
Documentation=https://github.com/baselithcore
# Start after networking (binds a host IP) and, if co-located, BaselithCore.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# ── EDIT THESE for the host ──────────────────────────────────────────────────
# The OS user that owns the checkout (NOT root).
User=CHANGE_ME
# Absolute path to the backstage-portal/ directory of the checkout.
WorkingDirectory=/opt/baselithcore/backstage-portal
# node 22 (or 24 — NOT 25) bin dir must be first on PATH. Point at your install:
#   nvm    → /home/<user>/.nvm/versions/node/v22.x.y/bin
#   distro → /usr/bin           brew → /home/linuxbrew/.linuxbrew/opt/node@22/bin
Environment=PATH=/home/CHANGE_ME/.nvm/versions/node/v22.23.1/bin:/usr/local/bin:/usr/bin:/bin
# ─────────────────────────────────────────────────────────────────────────────

# Host URLs + BaselithCore key live here (copied from backstage-portal.env.example).
EnvironmentFile=/etc/baselith/backstage-portal.env

# Run the vendored Yarn 4 release directly with node so PATH resolution can't
# pick a wrong yarn. `yarn start` == `backstage-cli repo start` (dev servers:
# frontend on APP_BASE_URL, backend on :7007).
ExecStart=/usr/bin/env node .yarn/releases/yarn-4.4.1.cjs start

Restart=on-failure
RestartSec=5
# backstage-cli responds to SIGINT for a clean shutdown of both dev servers.
KillSignal=SIGINT
TimeoutStopSec=20

# Light hardening (dev server needs write access to its own caches).
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
